# supervisior/supervisord.conf
[supervisord]
nodaemon=true

[program:mysqld]
command=/run-mysql.sh &

[program:phpfpm]
command=/usr/local/bin/run-phpfpm &

[program:nginx]
command=/usr/local/bin/start.sh &
